# configs/systemd/blind_control_hub.service
[Unit]
Description=Blind Control Hub
After=network.target
Wants=network.target

[Service]
Type=simple
User=sttark
Group=sttark
WorkingDirectory=/home/sttark/blind_control/hub
ExecStart=/usr/bin/python3 /home/sttark/blind_control/hub/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=PYTHONPATH=/home/sttark/blind_control

[Install]
WantedBy=multi-user.target
